# gateway/nidokey-gateway.service
# systemd unit del gateway. Instalar en /etc/systemd/system/nidokey-gateway.service
# Ajusta WorkingDirectory y User a tu despliegue. Luego:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now nidokey-gateway
#   sudo systemctl status nidokey-gateway

[Unit]
Description=Nidokey chat realtime gateway (WebSocket)
After=network.target

[Service]
Type=simple
# Usuario sin privilegios (créalo: sudo useradd --system --no-create-home nidokey)
User=nidokey
WorkingDirectory=/opt/nidokey-gateway
EnvironmentFile=/etc/nidokey-gateway.env
ExecStart=/usr/bin/node server.mjs
Restart=always
RestartSec=5
# Límites de escala: el default de fds (1024) topa en ~1k conexiones WS.
LimitNOFILE=65536
# Tope de RAM: si el proceso se descontrola, systemd lo mata y Restart lo relanza.
MemoryMax=4G
# Endurecimiento básico
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
